# ai-service/scripts/launchd/com.ringrift.owc-s3-mirror.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ringrift.owc-s3-mirror</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/armand/Development/RingRift/ai-service/scripts/owc_s3_mirror.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/armand/Development/RingRift/ai-service</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/armand/Development/RingRift/ai-service</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/armand</string>
    </dict>

    <!-- Run every 6 hours -->
    <key>StartInterval</key>
    <integer>21600</integer>

    <!-- Also run at load (startup) -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep alive - restart if it crashes -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/armand/Development/RingRift/ai-service/logs/owc_s3_mirror.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/armand/Development/RingRift/ai-service/logs/owc_s3_mirror.error.log</string>

    <!-- Nice level - run at lower priority -->
    <key>Nice</key>
    <integer>10</integer>

    <!-- Throttle interval - don't restart more than once per 5 minutes -->
    <key>ThrottleInterval</key>
    <integer>300</integer>
</dict>
</plist>
