# ai-service/deploy/systemd/elo-sync-coordinator.service
[Unit]
Description=RingRift Elo Database Sync Coordinator
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/ringrift/ai-service
ExecStart=/home/ubuntu/ringrift/ai-service/venv/bin/python scripts/elo_db_sync.py --mode coordinator --port 8766
Restart=always
RestartSec=10
Environment=PYTHONPATH=/home/ubuntu/ringrift/ai-service

[Install]
WantedBy=multi-user.target
