# ai-service/deploy/systemd/shadow-tournament.service
[Unit]
Description=RingRift Shadow Tournament Service
Documentation=https://github.com/ringrift/ai-service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/ringrift/ai-service

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/ubuntu/ringrift/ai-service"

# Main process
ExecStart=/home/ubuntu/ringrift/ai-service/venv/bin/python \
    scripts/shadow_tournament_service.py \
    --config config/unified_loop.yaml

# Restart policy
Restart=on-failure
RestartSec=60
StartLimitIntervalSec=600
StartLimitBurst=5

# Resource limits - tournaments are CPU intensive
MemoryMax=8G
CPUQuota=400%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shadow-tournament

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ReadWritePaths=/home/ubuntu/ringrift/ai-service/data
ReadWritePaths=/home/ubuntu/ringrift/ai-service/logs
ReadWritePaths=/home/ubuntu/ringrift/ai-service/models

[Install]
WantedBy=multi-user.target
