# ai-service/deploy/systemd/unified-data-sync.service
[Unit]
Description=RingRift Unified Data Sync Service
Documentation=https://github.com/ringrift/ai-service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/ringrift/ai-service

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/home/ubuntu/ringrift/ai-service"

# Main process - runs unified data sync with watchdog
# Note: unified_data_sync.py replaces the deprecated streaming_data_collector.py
ExecStart=/home/ubuntu/ringrift/ai-service/venv/bin/python \
    scripts/unified_data_sync.py \
    --config config/unified_loop.yaml \
    --hosts config/remote_hosts.yaml \
    --watchdog

# Restart policy
Restart=on-failure
RestartSec=15
StartLimitIntervalSec=300
StartLimitBurst=10

# Resource limits
MemoryMax=1G
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=unified-data-sync

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ReadWritePaths=/home/ubuntu/ringrift/ai-service/data
ReadWritePaths=/home/ubuntu/ringrift/ai-service/logs

[Install]
WantedBy=multi-user.target
